# atop.service
[Unit]
Description=advanced interactive monitor
After=syslog.target

[Service]
ExecStart=/usr/bin/atop

[Install]
WantedBy=multi-user.target
